# payload/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --noissue --skip-login --autologin min %I 38400 linux
